# deploy/sidecar/claude-tmux@.service
# deploy/sidecar/claude-tmux@.service
# Durable in-tmux swarph cell — systemd TEMPLATE (%i = cell name).
#
# One template replaces N hand-authored per-cell units. Enable a cell with:
#   systemctl --user enable --now claude-tmux@<cell>.service        (user cell)
#   sudo systemctl enable --now claude-tmux@<cell>.service          (root cell)
# Non-root user cells also need:  loginctl enable-linger $USER
#
# Shape below is the USER-unit form (swarph at ~/.local/bin, WantedBy=default.target).
# ROOT/system cell: swarph at /usr/local/bin, WantedBy=multi-user.target,
#   After=network-online.target. Adjust the two binary paths + the [Install] target.
#
# AI² fix (science-claude): KillMode=mixed + an explicit ExecStop kill-session —
#   the old `none` KillMode is deprecated (systemd warns). mixed = SIGTERM the main process,
#   SIGKILL the rest of the cgroup on stop; ExecStop owns the session teardown.
# BLOCKING fix (droplet): ExecStart is gated on `swarph cell verify %i` BEFORE the
#   per-name `has-session` check. The double-resume footgun is per-UUID, not per
#   tmux-name (two cell names pinning one UUID both pass has-session → both resume →
#   tmux-server death). verify refuses if the pinned UUID is already live under ANY
#   holder, and clears a stale poison-pin if the holder is dead.

[Unit]
Description=Durable in-tmux swarph cell %i — tmux session hosting a live claude
After=default.target

[Service]
Type=forking
Environment=TERM=xterm-256color
ExecStart=/bin/sh -lc '%h/.local/bin/swarph cell verify %i \
  && (/usr/bin/tmux has-session -t %i 2>/dev/null \
      || /usr/bin/tmux new-session -d -s %i "%h/.local/bin/swarph spawn %i")'
ExecStop=/usr/bin/tmux kill-session -t %i
KillMode=mixed
RemainAfterExit=yes

[Install]
WantedBy=default.target
